# vk-ads-telegram-bot.service
[Unit]
Description=VK Ads Manager Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/trouble/dev/vktest
Environment="PATH=/home/trouble/dev/vktest/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/trouble/dev/vktest/.venv/bin/python /home/trouble/dev/vktest/telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=append:/home/trouble/dev/vktest/telegram_bot.log
StandardError=append:/home/trouble/dev/vktest/telegram_bot.log

[Install]
WantedBy=multi-user.target
